# deploy/tradedash.service
# systemd unit for the venv deploy. Install to /etc/systemd/system/tradedashboard-backend.service
# Assumes: code at /home/ubuntu/app, venv at /home/ubuntu/app/backend/venv (python-dotenv
# loads /home/ubuntu/app/backend/.env itself — no EnvironmentFile needed).
#
#   sudo cp deploy/tradedash.service /etc/systemd/system/tradedashboard-backend.service
#   sudo systemctl daemon-reload && sudo systemctl enable --now tradedashboard-backend
#   journalctl -u tradedashboard-backend -f     # follow logs

[Unit]
Description=TradeDashBoard FastAPI backend
After=network.target network-online.target
Requires=network-online.target

[Service]
Type=simple
User=ubuntu
Group=ubuntu
WorkingDirectory=/home/ubuntu/app/backend
ExecStart=/home/ubuntu/app/backend/venv/bin/python -u /home/ubuntu/app/backend/run.py
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal
# One instance only — FYERS allows a single websocket connection per app.

[Install]
WantedBy=multi-user.target
